# crontab
SPARQL_ENDPOINT=$SPARQL_ENDPOINT
DATA=$DATA

0 * * * * /app/crontask.sh >> /var/log/cron.log 2>&1
